# com.screenshot.uploader.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>EnvironmentVariables</key>
	<dict>
		<key>PATH</key>
		<string>/Users/m5/screenshot/.venv/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
	</dict>
	<key>KeepAlive</key>
	<true/>
	<key>Label</key>
	<string>com.screenshot.uploader</string>
	<key>ProgramArguments</key>
	<array>
		<string>/Users/m5/screenshot/.venv/bin/python</string>
		<string>/Users/m5/screenshot/screenshot_uploader.py</string>
		<string>--path</string>
		<string>/Users/m5/screenshot/temp</string>
		<string>--limit</string>
		<string>20</string>
	</array>
	<key>RunAtLoad</key>
	<true/>
	<key>StandardErrorPath</key>
	<string>/Users/m5/screenshot/screenshot_uploader.err</string>
	<key>StandardOutPath</key>
	<string>/Users/m5/screenshot/screenshot_uploader.out</string>
	<key>WorkingDirectory</key>
	<string>/Users/m5/screenshot</string>
</dict>
</plist>
