# com.user.screenshot_uploader.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.user.screenshot_uploader</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/m5/screenshot/.venv/bin/python</string>
        <string>/Users/m5/screenshot/screenshot_uploader.py</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/m5/screenshot</string>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardErrorPath</key>
    <string>/Users/m5/screenshot/screenshot_uploader.err</string>
    <key>StandardOutPath</key>
    <string>/Users/m5/screenshot/screenshot_uploader.out</string>
</dict>
</plist> 
